# deployment/trading-engine.service
[Unit]
Description=Trading Engine Service
After=network.target mongodb.service
Wants=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/Independent trade engine
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /path/to/Independent trade engine/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Optional: Resource limits
# MemoryLimit=2G
# CPUQuota=100%

# Optional: Security hardening
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ReadWritePaths=/path/to/Independent trade engine

[Install]
WantedBy=multi-user.target
